# com.tiktok.autopost.evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tiktok.autopost.evening</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Library/WebServer/Documents/TIKTOK/auto-post.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Library/WebServer/Documents/TIKTOK</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Library/WebServer/Documents/TIKTOK/logs/autopost-evening.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/WebServer/Documents/TIKTOK/logs/autopost-evening-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/masscom/.nvm/versions/node/v24.14.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
